# linux/systemd/user/xremap.service
[Unit]
Description=Start remap

[Service]
Type=exec
TimeOutStartSec=30
WorkingDirectory=%h
StandardOutput=journal
ExecStart=xremap --watch=device %h/.config/xremap/config.yml
Restart=always

[Install]
WantedBy=multi-user.target
